# packages/linuxadmin/linuxadmin-0.0.1.tar.gz/linuxadmin-0.0.1/linuxadmin.service
[Unit]
Description=LinuxAdmin thinclient
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/linuxadmin
Restart=on-failure

[Install]
WantedBy=multi-user.target