# packages/criscostack-brik/criscostack_brik-1.0.0-py3-none-any.whl/brik/config/templates/systemd/frappe-bench-redis.target
[Unit]
After=network.target
Wants={{ brik_name }}-redis-cache.service {{ brik_name }}-redis-queue.service {{ brik_name }}-redis-socketio.service

[Install]
WantedBy=multi-user.target
